# etc/gem.service
[Unit]
Description=GEM Motion Control Daemon

[Service]
ExecStart=/usr/local/sbin/gem-motiond
Restart=always

[Install]
WantedBy=multi-user.target
